# roles/groovy/files/etc/profile.d/groovy.sh
export PATH=$PATH:/usr/local/groovy/bin
