# roles/gradle/files/etc/profile.d/gradle.sh
export PATH=$PATH:/usr/local/gradle/bin
